# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:postgres]
command=/bin/bash /postgis_run.sh

[program:tomcat]
command=/bin/bash /tomcat_run.sh
